# config/ghost-config/ghost.conf
# ghost

start on startup

script
  cd /srv/www/ghost
  NODE_ENV=production npm start
end script