# config/ghost-config/ghost.conf
# ghost

start on startup
stop on shutdown
respawn
respawn limit 99 5

script
  cd /srv/www/ghost
  NODE_ENV=production npm start --production
end script
